# configs/etc/s6/insync/finish
#!/usr/bin/env bash

/usr/bin/insync-headless quit
